# sv/firefox-overlay-cron/run
#!/bin/sh
exec snooze -H \* -- firefox-overlay flush
